# fiberglass.service
[Unit]
Description=Fiberglass Processing Insulation

[Service]
ExecStart=start-container
Restart=always

[Install]
WantedBy=default.target
